# systemd/tl-events.timer
[Unit]
Description=Run tl-events.service every 41 minutes

[Timer]
OnCalendar=*:0/41

[Install]
WantedBy=timers.target