# systemd/tl-scraper.timer
[Unit]
Description=Run tl-scraper.service every 10 minutes

[Timer]
OnCalendar=*:0/59
