# units/calendar.timer
[Unit]
Description=Run calendar.service every 95 minutes

[Timer]
OnCalendar=*0:0:0/95

[Install]
WantedBy=timers.target