# units/calendar.service.in
[Unit]
Description=Scraper for ESports calendar
After=network.target

[Service]
WorkingDirectory=BIN_DIR
Type=oneshot

ExecStart=BIN_DIR/calendar --calendar tl hs dota lol pfw

