# units/calendarevents.service.in
[Unit]
Description=Scraper for ESports events
After=network.target

[Service]
WorkingDirectory=BIN_DIR
Type=oneshot

ExecStart=BIN_DIR/calendarevents --calendar sc2 dota lol csgo

