# units/calendarevents.timer
[Unit]
Description=Run tl-events.service every 211 minutes

[Timer]
OnCalendar=*-*-* 00:00:00/211

[Install]
WantedBy=timers.target