# units/calendar.timer
[Unit]
Description=Run calendar.service every 95 minutes

[Timer]
OnCalendar=*-*-* 00:00:00/95

[Install]
WantedBy=timers.target