# units/calendarevents.service.in
[Unit]
Description=Event scraper for ESports events
After=network.target

[Service]
WorkingDirectory=BIN_DIR
Type=oneshot

ExecStart=BIN_DIR/calendarevents --calendar tl hs dota lol pfw


