# units/calendar.service.in
[Unit]
Description=Calendar scraper for ESports events
After=network.target

[Service]
WorkingDirectory=BIN_DIR
Type=oneshot

ExecStart=BIN_DIR/calendar --calendar sc2 brw csg hot sma dot lol hrt pfw

