# units/calendarevents.timer
[Unit]
Description=Run tl-events.service every 211 minutes

[Timer]
OnCalendar=*0:0:0/211

[Install]
WantedBy=timers.target